# openwrt/package/hexadaemon/hexadaemon.init
#!/bin/sh /etc/rc.common
# Copyright (c) 2013 Stephan Platz

START=94

interface=`uci get network.wan.ifname`

start()
{
	ADDR="";
	for x in `ip -6 a show scope global | sed "s/.*inet6 \([^ \/]*\).*/\1/;tx;d;:x"`;
		do ADDR="$ADDR-a $x ";
	done;
	/usr/bin/hexadaemon -i 2 -I ${interface} ${ADDR}
}

stop()
{
	killall hexadaemon
}
