# openwrt/package/hexadaemon/hexadaemon.init
#!/bin/sh /etc/rc.common
# Copyright (c) 2013 Stephan Platz

START=94

interface=`uci get network.wan.ifname`

start()
{
	/usr/bin/hexadaemon -i 2 -I ${interface}
}

stop()
{
	killall hexadaemon
}
