# openwrt/package/luci/applications/luci-flukso/root/etc/uci-defaults/luci-flukso
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@flukso[-1]
	add ucitrack flukso
	set ucitrack.@flukso[-1].init=msg-flukso2
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
